# profile.d/ls.sh
alias ls="ls -aFh --color=auto"

eval `dircolors $HOME/.system/files/dircolors.ansi-dark`
